# rtsp-snapshot.timer
[Unit]
Description=Timer for RTSP snapshot
[Timer]
OnCalendar=*-*-* 5:0:0
AccuracySec=10s

[Install]
WantedBy=timers.target
